# deploy/systemd/schule-postgres-backup.service.example
[Unit]
Description=PostgreSQL-Backup der schulischen Zeugnisverwaltung
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
User=__SERVICE_USER__
Group=__SERVICE_GROUP__
WorkingDirectory=__PROJECT_DIR__
Environment=BACKUP_DIR=__BACKUP_DIR__
Environment=BACKUP_RETENTION_DAYS=30
ExecStart=__PROJECT_DIR__/scripts/backup-postgres.sh
UMask=0077
Nice=10
NoNewPrivileges=true
PrivateTmp=true
