# deploy/systemd/schule-compose.service.example
[Unit]
Description=Schulische Zeugnisverwaltung per Docker Compose
Requires=docker.service
After=docker.service network-online.target

[Service]
Type=oneshot
WorkingDirectory=__PROJECT_DIR__
ExecStart=/usr/bin/docker compose up -d --remove-orphans
ExecStop=/usr/bin/docker compose stop
RemainAfterExit=yes
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
